# infra/vespergrid-api.service
[Unit]
Description=VesperGrid API (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/vespergrid/api.env
WorkingDirectory=/opt/vespergrid
ExecStart=/opt/vespergrid-venv/bin/python -m uvicorn vespergrid.main:app --host 0.0.0.0 --port 8742
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
